# assets/launchd.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wx-article-sync.daily</string>

  <key>WorkingDirectory</key>
  <string>/path/to/wx-article-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/absolute/path/to/uv</string>
    <string>run</string>
    <string>wx-article-sync</string>
    <string>--config</string>
    <string>config.json</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/path/to/logs/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/logs/launchd.stderr.log</string>
</dict>
</plist>

